# supervisor/guacd.sv.conf
[program:guacd]
command=guacd -f -C /etc/ssl/certs/guacd.pem -K /etc/ssl/private/guacd.key
startretries=0
autorestart=true
redirect_stderr=true